# scripts/pypilot-route-nav.service
[Unit]
Description=Pypilot Route Navigation Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pypilot
Environment=PYTHONUNBUFFERED=1
Environment=PYPILOT_ROUTE_SERVICE_PORT=20221
ExecStart=/usr/bin/python3 -m web.route_service
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
